# supervisord.conf
[supervisord]
nodaemon=true

[program:octane]
process_name=%(program_name)s
command=php /app/artisan octane:start --server=frankenphp
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/octane.log
stderr_logfile=/var/log/octane.error.log